# deploy/systemd/campus-course.service
[Unit]
Description=Campus Assistant Course Service
After=network.target campus-gateway.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/campus
ExecStart=/usr/bin/java -jar /opt/campus/campus-assistant/course-service/target/course-service-1.0-SNAPSHOT.jar
Restart=always
RestartSec=5
StandardOutput=append:/opt/campus/logs/course-service.log
StandardError=append:/opt/campus/logs/course-service.log

[Install]
WantedBy=multi-user.target
